# com.adnsys.robo-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adnsys.robo-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>PASTA_RAIZ/watchdog.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>PASTA_RAIZ</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>PASTA_RAIZ/watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>PASTA_RAIZ/watchdog.log</string>
</dict>
</plist>
